# src/server/10.10/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.server.eventsd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.xserver.eventsd</key>
		<true/>
	</dict>
	<key>Label</key>
	<string>com.apple.server.eventsd</string>
	<key>KeepAlive</key>
	<true/>
	<key>OnDemand</key>
	<false/>
	<key>Disabled</key>
	<false/>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Server.app/Contents/ServerRoot/usr/libexec/ServerEventsDaemon</string>
	</array>
</dict>
</plist>
