# src/server/10.10/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.afctl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.afctl</string>
	<key>Disabled</key>
	<true/>
	<key>Program</key>
	<string>/Applications/Server.app/Contents/ServerRoot/usr/libexec/afctl</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Server.app/Contents/ServerRoot/usr/libexec/afctl</string>
	</array>
	<key>StartInterval</key>
	<integer>900</integer>
	<key>runAtLoad</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
</dict>
</plist>
