# src/server/10.13/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.push_notify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>Label</key>
	<string>com.apple.push_notify</string>
	<key>KeepAlive</key>
	<true/>
	<key>PreventsSleep</key>
	<true/>
	<key>Program</key>
	<string>/Applications/Server.app/Contents/ServerRoot/usr/sbin/push_notify</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Server.app/Contents/ServerRoot/usr/sbin/push_notify</string>
	</array>
</dict>
</plist>
