# src/server/10.13/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.server.httpd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>Label</key>
	<string>com.apple.server.httpd</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>SERVER_INSTALL_PATH_PREFIX</key>
		<string>/Applications/Server.app/Contents/ServerRoot</string>
		<key>MODULE_INSTALL_PATH_PREFIX</key>
		<string/>
		<key>XPC_SERVICES_UNAVAILABLE</key>
		<string>1</string>
	</dict>
	<key>PreventsSleep</key>
	<true/>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/httpd-wrapper</string>
                <string>-D</string>
		<string>FOREGROUND</string>
		<string>-f</string>
		<string>/Library/Server/Web/Config/apache2/httpd_server_app.conf</string>
	        <string>-E</string>
                <string>/var/log/apache2/error_log</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>ExitTimeOut</key>
	<integer>61</integer>
</dict>
</plist>
