# src/server/10.14/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.server.alertsd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.server.alertsd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.xserver.alerts</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Server.app/Contents/ServerRoot/usr/libexec/AlertsDaemon</string>
	</array>
	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
