# src/server/10.11/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.disks.smart.status.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<false/>
	<key>Label</key>
	<string>com.apple.disks.smart.status</string>
	<key>Disabled</key>
	<false/>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Server.app/Contents/ServerRoot/usr/libexec/smartmonitor.sh</string>
	</array>
	<key>StartOnMount</key>
	<true/>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
	</dict>
</dict>
</plist>
