# src/server/10.11/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.serverdocs.sdhd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>com.apple.serverdocs.sdhd</string>
        <key>POSIXSpawnType</key>
        <string>Interactive</string>
        <key>MultipleInstances</key>
        <true/>
        <key>ProgramArguments</key>
        <array>
                <string>/Applications/Server.app/Contents/ServerRoot/System/Library/PrivateFrameworks/ServerDocsHelper.framework/sdhd</string>
        </array>
        <key>MachServices</key>
        <dict>
                <key>com.apple.serverdocs.sdhd</key>
                <true/>
        </dict>
        <key>EnvironmentVariables</key>
        <dict>
                <key>SDHD_SANDBOX</key>
                <string>/Applications/Server.app/Contents/ServerRoot/System/Library/Sandbox/com.apple.serverdocs.sdhd.sb</string>
                <key>SDH_FRAMEWORK_DIR</key>
                <string>/Applications/Server.app/Contents/ServerRoot/System/Library/PrivateFrameworks/ServerDocsHelper.framework</string>
        </dict>
</dict>
</plist>
