# src/server/10.9/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.collabd.stats.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>Label</key>
	<string>com.apple.collabd.stats</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Server.app/Contents/ServerRoot/usr/bin/wikiadmin</string>
		<string>gatherWikiStats</string>
	</array>
	<key>UserName</key>
	<string>_teamsserver</string>
	<key>GroupName</key>
	<string>_teamsserver</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
</dict>
</plist>
