# src/server/10.9/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.servermgr_xcode.get_currently_connected_devices.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.servermgr_xcode.get_currently_connected_devices</string>
	<key>Disabled</key>
	<false/>
	<key>PreventsSleep</key>
	<true/>
	<key>UserName</key>
	<string>root</string>
	<key>GroupName</key>
	<string>wheel</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Server.app/Contents/ServerRoot/usr/sbin/serveradmin</string>
		<string>command</string>
		<string>xcode:command= getCurrentlyConnectedDevices</string>
	</array>
	<key>WatchPaths</key>
	<array>
		<string>/Library/Server/Xcode/Cache/KickFiles/com.apple.servermgr_xcode.get_currently_connected_devices/</string>
	</array>
	<key>StartInterval</key>
	<integer>600</integer>
</dict>
</plist>
