# src/server/10.9/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.salearn.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.salearn</string>
	<key>PreventsSleep</key>
	<true/>
	<key>Program</key>
	<string>/Applications/Server.app/Contents/ServerRoot/usr/libexec/spamassassin/sa_learn.pl</string>
	<key>ProgramArguments</key>
	<array>
		<string>sa_learn.pl</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>2</integer>
		<key>Minute</key>
		<integer>15</integer>
	</dict>
</dict>
</plist>
