# src/server/10.9/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.mail_status_init.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.mail_status_init</string>
	<key>Program</key>
	<string>/Applications/Server.app/Contents/ServerRoot/usr/sbin/serveradmin</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Server.app/Contents/ServerRoot/usr/sbin/serveradmin</string>
		<string>command</string>
		<string>mail:command</string>
		<string>=</string>
		<string>doStatusCheck</string>
	</array>
	<key>UserName</key>
	<string>root</string>
	<key>StartInterval</key>
	<integer>7200</integer>
</dict>
</plist>
