# src/server/10.9/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.xcscredhandler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.xcscredhandler</string>
	<key>Disabled</key>
	<false/>
	<key>OnDemand</key>
	<false/>
	<key>UserName</key>
	<string>_xcscredserver</string>
	<key>GroupName</key>
	<string>_xcscredserver</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Server.app/Contents/ServerRoot/System/Library/CoreServices/XCSCredentialServer.bundle/Contents/MacOS/xcscredhandler</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.xcscredhandler</key>
		<true/>
	</dict>
</dict>
</plist>
