# src/server/10.9/Payload/Applications/Server.app/Contents/ServerRoot/private/etc/jabberd/jabberd-s2s.conf
description	"jabberd Server-2-Server component"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec sudo -u jabber /Applications/Server.app/Contents/ServerRoot/usr/libexec/jabberd/s2s
