# src/server/10.12/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.DeviceManagement.DMHelper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>BinaryOrderPreference</key>
  <array>
    <integer>7</integer>
    <integer>-1</integer>
  </array>
  <key>EnableTransactions</key>
  <true/>
  <key>EnablePressuredExit</key>
  <true/>
  <key>Label</key>
  <string>com.apple.DeviceManagement.DMHelper</string>
  <key>MachServices</key>
  <dict>
    <key>com.apple.DeviceManagement.DMHelper</key>
    <true/>
  </dict>
  <key>ProcessType</key>
  <string>Adaptive</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Server.app/Contents/ServerRoot/usr/share/devicemgr/dm_helper</string>
  </array>
  <key>StandardErrorPath</key>
  <string>/dev/null</string>
  <key>StandardOutPath</key>
  <string>/dev/null</string>
</dict>
</plist>
