# src/server/10.12/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.serverdocs.authserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>com.apple.serverdocs.authserver</string>
        <key>RunAtLoad</key>
        <true/>
        <key>KeepAlive</key>
        <true/>
        <key>POSIXSpawnType</key>
        <string>Interactive</string>
        <key>EnablePressuredExit</key>
        <false/>
        <key>ProgramArguments</key>
        <array>
                <string>/Applications/Server.app/Contents/ServerRoot/usr/libexec/authserver</string>
        </array>
        <key>MachServices</key>
        <dict>
                <key>com.apple.serverdocs.authserver</key>
                <true/>
        </dict>
</dict>
</plist>
