# src/server/10.12/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.DeviceManagement.dmrunnerd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>EnableTransactions</key>
  <false/>
  <key>GroupName</key>
  <string>_devicemgr</string>
  <key>KeepAlive</key>
  <true/>
  <key>Label</key>
  <string>com.apple.DeviceManagement.dmrunnerd</string>
  <key>PreventsSleep</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Server.app/Contents/ServerRoot/usr/share/devicemgr/dmrunnerd</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>UserName</key>
  <string>_devicemgr</string>
</dict>
</plist>
