# src/server/10.12/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.collabd.expire.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>Label</key>
	<string>com.apple.collabd.expire</string>
	<key>PreventsSleep</key>
	<true/>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Server.app/Contents/ServerRoot/usr/share/collabd/server/ruby/expire-collabd-data</string>
	</array>
	<key>StartInterval</key>
	<integer>600</integer>
	<key>UserName</key>
	<string>_teamsserver</string>
	<key>GroupName</key>
	<string>_teamsserver</string>
</dict>
</plist>
