# src/server/10.12/Payload/Applications/Server.app/Contents/ServerRoot/System/Library/LaunchDaemons/com.apple.serverdocs.quicklook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.serverdocs.quicklook</string>
	<key>PreventsSleep</key>
	<true/>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/QuickLook.framework/Resources/quicklookd.app/Contents/MacOS/quicklookd</string>
		<string>--servername</string>
		<string>com.apple.serverdocs.quicklook</string>
		<string>--allowedplugins</string>
		<string>/Library/Server/ServerDocs/Config/qlgenerators.plist</string>
		<string>--nouser</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.serverdocs.quicklook</key>
		<true/>
	</dict>
	<key>ThrottleInterval</key>
	<integer>1</integer>
</dict>
</plist>
